# etc/init.d/cg_bg
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description='setup cgroupv2 "background"'
command=''

depend() {
    need cgroups
}


cgroup2_remove() {
	einfo skipping cgroup2 remove
}

start() {
	CG_BG_DEFAULT_IOLAT=1000000
	CG_BG_EXTRA_DEFAULT_DEVICES="sda sdb sbc sdd"
    cgroup_path="$(cgroup2_find_path)"
	[ -n "${cgroup_path}" ] || return 1
	[ -d "${cgroup_path}" ] || mkdir "${cgroup_path}"
    rc_cgroup_path="${cgroup_path}/${RC_SVCNAME}"
	[ -d "${rc_cgroup_path}" ] || mkdir "${rc_cgroup_path}"
	
	#echo $(gen_io_lat_strings ${CG_BG_DEFAULT_IOLAT} ${CG_BG_EXTRA_DEFAULT_DEVICES}) > "${rc_cgroup_path}/io.latency"
	#echo $(gen_io_lat_strings ${CG_BG_DEFAULT_IOLAT} ${CG_BG_EXTRA_DEFAULT_DEVICES}) \> "${rc_cgroup_path}/io.latency"

	#lat_strings=($(gen_io_lat_strings $CG_BG_DEFAULT_IOLAT $CG_BG_EXTRA_DEFAULT_DEVICES))

	#IFS="\n"
	#for line in $lat_strings; do
	#	
	#	einfo "\"${line}\"" \> "${rc_cgroup_path}/io.latency"
	#	
	#	echo "${line}" > "${rc_cgroup_path}/io.latency"
	#done
}

stop(){
	einfo ""
}
