# etc/init.d/cluster-service
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="${RC_SVCNAME}"

description="cluster thing ${name}"

depend() {
if [[ "$RC_SVCNAME" != cluster-service ]];then
	einfo "debug: cluster service $name needs "cluster_${providingnode}_${name}""
	need "cluster_${providingnode}_${name}"
    need cluster-com
fi
}

start(){
	ebegin "starting ${name} .."
}

stop(){
    rc-service "remote_${providingnode}_${name}" stop
	eend
}












