# etc/init.d/in_tftpd_node
#!/sbin/openrc-run
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

. /lib/rc/sh/rc-linkname-parsers.sh
parse_svcdata
providingnode="${svcdata[node]}"
servicename="${svcdata[name]}"

depend() {
    if [[ "${RC_SVCNAME}" != in_tftp_node ]];then
        need "net.${RC_SVCNAME}"
    fi
}

start() {
	ebegin "Starting bootdata tftpd for ${INTFTPD_NODENAME}"
	einfo "addr=${INTFTPD_ADDRESS} path=${INTFTPD_PATH} nodename=${INTFTPD_NODENAME}"
	start-stop-daemon \
	--start --exec /usr/sbin/in_tftpd_${INTFTPD_NODENAME} \
	-- \
    --user "in_tftpd_${INTFTPD_NODENAME}" \
    --listen \
    -vvv --ipv4 \
    --permissive \
    --port-range "24510:24550" \
    --pidfile "/var/run/in_tftpd_${INTFTPD_NODENAME}" \
	--address "${INTFTPD_ADDRESS}" \
	--secure "${INTFTPD_PATH}"
	eend $?
}

stop() {
	ebegin "Stopping bootdata tftpd for ${INTFTPD_NODENAME}"
	start-stop-daemon \
    --pidfile "/var/run/in_tftpd_${INTFTPD_NODENAME}" \
	--name "/usr/sbin/in_tftpd_${INTFTPD_NODENAME}" \
	--stop
	eend $?
}
