# etc/init.d/link-name-provide
#!/sbin/openrc-run

depend(){
 	provide "${RC_SVCNAME}"
}

start(){
	einfo "Service ${RC_SVCNAME} is provided."
}
