# etc/init.d/nbd-client
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

# the client reads the "/etc/nbdtab" file

name="nbd-client daemon"
description=""

command=nbd-client

devices=$(cat /etc/nbdtab|grep -Ev '^\W*\#'|cut -d' ' -f1)

pidfile="/var/run/$(hostname)/${RC_SVCNAME}.pid"

#start_stop_daemon_args=""
#start_stop_daemon_args="${start_stop_daemon_args} --stdout-logger logger\ -p\ daemon.info"
#start_stop_daemon_args="${start_stop_daemon_args} --stderr-logger logger\ -p\ daemon.crit"
#start_stop_daemon_args="${start_stop_daemon_args} --ionice 3"
#start_stop_daemon_args="${start_stop_daemon_args} --pidfile ${pidfile} --make-pidfile"
#start_stop_daemon_args="${start_stop_daemon_args} --background"
#start_stop_daemon_args="${start_stop_daemon_args} --wait 10000"

depend()
{
    need net.dusteater iptables ip6tables
}

start()
{
	for d in $devices;do
		nbd_free_dev ${d}
		command_args=$i #todo check
		nbd-client $d
        pid=$(nbd-client -c /dev/${d})
    	ret=$?
    	if [[ $ret -eq 0 ]];then
        	echo ${pid} >> ${pidfile}
    	else
        	eend 1 "failed to start ${RC_SVCNAME}"
    	fi
	done
}

stop()
{	for d in $devices;do
    	nbd_free_dev ${d}
		pfc=$(cat ${pidfile}|tail -n+1)
    	default_stop
		echo ${pfc} > ${pidfile}
	done
}

nbd_free_dev()
{
    einfo freeing device
    nbd-client -d /dev/${1}
    [[ $? -ne 2 ]] && eend 0 "/dev/${d} free"
    einfo "device ${d} clear"
}
