# etc/init.d/nbd-server-tmp
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

# the server reads the "/etc/nbd-server/config" file

name="nbd-server daemon tmp"
description="the part of the nbd-server that provides the tmp export"

depend() {
need nbd-server
}

start(){
	einfo "starting nbd-server-tmp"
}
