# etc/init.d/nbd-swap
#!/sbin/openrc-run
# Copyright (c) 2007-2015 The OpenRC Authors.
# See the Authors file at the top-level directory of this distribution and
# https://github.com/OpenRC/openrc/blob/master/AUTHORS
#
# This file is part of OpenRC. It is subject to the license terms in
# the LICENSE file found in the top-level directory of this
# distribution and at https://github.com/OpenRC/openrc/blob/master/LICENSE
# This file may not be copied, modified, propagated, or distributed
# except according to the terms contained in the LICENSE file.

name="nbd-swap"
description="mounts tmpfs over nfs , creates swap file and mount that via \"nbd\""

pidfile_dir="/var/run/$(hostname)"
pidfile="${pidfile_dir}/${RC_SVCNAME}.pid"

depend()
{
    need dusteater net.dusteater iptables ip6tables
	after clock root
	keyword -docker -jail -lxc -openvz -prefix -systemd-nspawn -vserver
}

start()
{
	        ebegin "Activating ${RC_SVCNAME} ... "
		    mount /mnt/dusteater/var/tmp/esadc \
            && free_space_k="$(df -h -Bk --output=avail /mnt/dusteater/var/tmp/esadc/|tail -n+2|head -n1|sed 's/[^0-9]//g')" \
            && allocate_space_k=$(( $free_space_k - 1000 )) \
            &&[[ $allocate_space_k -gt 0 ]] \
		   && fallocate --length "${allocate_space_k}k" /mnt/dusteater/var/tmp/esadc/esadc_swapfile \
		   && chmod 0600 /mnt/dusteater/var/tmp/esadc/esadc_swapfile \
           && modprobe nbd \
           && sudo ssh root@de 'rc-service nbd-server start' \
           && start_nbd_client \
		   && mkswap /dev/nbd0 \
		   && swapon /dev/nbd0 \
            && cat /proc/sys/vm/overcommit_ratio > /var/run/sysctl_backup.vm.overcommit_ratio \
           && echo 50 > /proc/sys/vm/overcommit_ratio \
            && cat /proc/sys/vm/watermark_scale_factor > /var/run/sysctl_backup.vm.watermark_scale_factor \
           && echo 1000 > /proc/sys/vm/watermark_scale_factor \
            && cat /proc/sys/vm/swappiness > /var/run/sysctl_backup.vm.swappiness \
           && echo 100 > /proc/sys/vm/swappiness \
           && echo -n
	# eend 0 # If swapon has nothing todo it errors, so always return 0
    # but should have something to do
}

stop()
{
	ebegin "Deactivating ${RC_SVCNAME} ..."
		     	swapoff /dev/nbd0 \
                && cat /var/run/sysctl_backup.vm.overcommit_ratio > /proc/sys/vm/overcommit_ratio \
                && cat /var/run/sysctl_backup.vm.watermark_scale_factor > /proc/sys/vm/watermark_scale_factor \
                && cat /var/run/sysctl_backup.vm.swappiness > /proc/sys/vm/swappiness \
                && nbd_free_dev nbd0 \
                && sudo ssh root@de rc-service nbd-server stop \
                && rm -f /mnt/dusteater/var/tmp/esadc/esadc_swapfile \
				&& umount /mnt/dusteater/var/tmp/esadc \
                && ps --no-header -p $(cat "${pidfile}"|head -n1) && eend 1 "it seems that nbd client did not exit properly"
                pfc=$(cat ${pidfile}|tail -n+1) \
                &&echo ${pfc} > ${pidfile}
	eend 0
}

start_nbd_client()
{
		nbd_free_dev nbd0
        [[ -d  $(realpath "${pidfile_dir}") ]] || mkdir -p "${pidfile_dir}"
        chown -R nbd:nbd "${pidfile_dir}"
        chmod -R ug=rwX,o= "${pidfile_dir}"
        chmod g+s "${pidfile_dir}"
		nbd-client -name swap -swap -persist -connections 8 dusteater 8464 /dev/nbd0
        pid=$(nbd-client -c /dev/nbd0)
    	ret=$?
    	if [[ $ret -eq 0 ]];then
            echo ${pid} >> ${pidfile}
    	else
        	eend 1 "failed to start ${RC_SVCNAME}"
    	fi
}

nbd_free_dev()
{
    einfo freeing device
    nbd-client -d /dev/${1}
    [[ $? -ne 2 ]] && eend 0 "/dev/${1} free"
    einfo "device ${1} clear"
}
