# etc/init.d/nbd-swap
#!/sbin/openrc-run
# Copyright (c) 2007-2015 The OpenRC Authors.
# See the Authors file at the top-level directory of this distribution and
# https://github.com/OpenRC/openrc/blob/master/AUTHORS
#
# This file is part of OpenRC. It is subject to the license terms in
# the LICENSE file found in the top-level directory of this
# distribution and at https://github.com/OpenRC/openrc/blob/master/LICENSE
# This file may not be copied, modified, propagated, or distributed
# except according to the terms contained in the LICENSE file.

name="nbd-swap"
description="enables network block device based swap"

nbd_dev_link_name="${RC_SVCNAME##*_}"
dev_path="/dev/${nbd_dev_link_name}"

depend()
{
    need "${nbd_dev_link_name}"
}

start()
{
	        ebegin "Activating ${RC_SVCNAME}"
	        einfo "${dev_path}"
			mkswap "${dev_path}" || eend 1 || return 
		    swapon "${dev_path}" || eend 1 || return

            cat /proc/sys/vm/overcommit_ratio > /var/run/sysctl_backup.vm.overcommit_ratio || ewarn "${RC_SVCNAME}"
            cat /proc/sys/vm/watermark_scale_factor > /var/run/sysctl_backup.vm.watermark_scale_factor || ewarn "${RC_SVCNAME}"
            cat /proc/sys/vm/swappiness > /var/run/sysctl_backup.vm.swappiness || ewarn "${RC_SVCNAME}"

            echo 50 > /proc/sys/vm/overcommit_ratio || ewarn "${RC_SVCNAME}"
            echo 1000 > /proc/sys/vm/watermark_scale_factor || ewarn "${RC_SVCNAME}"
            echo 100 > /proc/sys/vm/swappiness || ewarn "${RC_SVCNAME}"

			chmod 0700 "${dev_path}" || eerror "${RC_SVCNAME}" 
            eend 0
}

stop()
{
	ebegin "Deactivating ${RC_SVCNAME}"
	einfo "${dev_path}"
	swapoff "${dev_path}" || eend 1 || return
    cat /var/run/sysctl_backup.vm.overcommit_ratio > /proc/sys/vm/overcommit_ratio || ewarn "${RC_SVCNAME}"
    cat /var/run/sysctl_backup.vm.watermark_scale_factor > /proc/sys/vm/watermark_scale_factor || ewarn "${RC_SVCNAME}"
    cat /var/run/sysctl_backup.vm.swappiness > /proc/sys/vm/swappiness || ewarn "${RC_SVCNAME}"
	eend 0
}

