# etc/init.d/nbd-server-part
#!/sbin/openrc-run
# Copyright 2019 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

# the server reads the "/etc/nbd-server/config" file

part_name="${RC_SVCNAME#*_}"

name="${RC_SVCNAME}"
description="The part of the nbd-server that provides the \"${part_name}\" export."


depend() {
		need nbd-server
}

start(){
	einfo "starting ${RC_SVCNAME}"
}

stop(){
	einfo "stopping ${RC_SVCNAME}"
}
