# etc/init.d/nbd-fs-backing-tmp
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="nbd-fs-backing-tmp"
description="enables the network block device based backing device for the filesystem at /tmp"

depend() {
	need nbd-ramdisk-tmp
}

start() {
 ebegin "adding dev \"/dev/nbd-ramdisk-tmp\" to btrfs-fs at \"/run/tmp\""
 if [[ $(df -h /run/tmp | tail -n+2|cut -d' ' -f1) = "/dev/loop0" ]];then
 	if [[ $(losetup -O back-file /dev/loop0 |tail -n+2) = /run/tmp.img ]];then
		btrfs dev add /dev/nbd-ramdisk-tmp /run/tmp
		eend 0 ;return
	fi
 fi
 eend 1 return
}

stop() {
 ebegin "removing dev \"/dev/nbd-ramdisk-tmp\" from btrfs-fs at \"/run/tmp\""
 if [[ $(df -h /run/tmp | tail -n+2|cut -d' ' -f1) = "/dev/loop0" ]];then
 	if [[ $(losetup -O back-file /dev/loop0 |tail -n+2) = /run/tmp.img ]];then
		btrfs dev del /dev/nbd-ramdisk-tmp /run/tmp
		eend 0 ;return
	fi
 fi
 eend 1 return
}
