# etc/init.d/net.node
#!/sbin/openrc-run

. /lib/rc/sh/rc-linkname-parsers.sh
parse_svcdata
node="${svcdata[node]}"
ifname="${svcdata[name]}"

depend() {
if [[ "$RC_SVCNAME" != net.node ]];then
        need "net.${ifname}"
        provide "net.${node}"
fi
}

start(){
	einfo "Starting network interface for node \"${node}\"."
}
