# etc/init.d/node-bootsrv
#!/sbin/openrc-run

. /lib/rc/sh/rc-linkname-parsers.sh
parse_svcdata
node="${svcdata[node]}"
thing="${svcdata[name]}"

depend() {
if [[ "$RC_SVCNAME" != node-bootsrv ]];then
        need "${thing}"
        provide "node-bootsrv_${node}"
fi
}

start(){
	einfo "Starting boot server for node \"${node}\"."
}
