# etc/init.d/run_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="run_subdirs"
description="population subdirs and stuff inside of /run"

# Specify a location of a directory tree
# with the following config variable

# run_tree="${run_tree}"

# the contents of the tree will be copied
# to /run with \"cp -a --no-clobber\"


depend() {
need root
}

start() {
    einfo "populating run subdirs"
    for node in $(cat /etc/nodelist);do
        einfo "node ${node}"

        if [[ "${node}" != "$(hostname)" ]];then
            mkdir -p   "/run/${node}"
            chmod 0751 "/run/${node}"
        else
            # ln -s "." "/run/${node}"
            linkname="/run/${node}"
            linktarget="."
            if ! [[ -e "${linkname}" ]];then
                ln -s  "${linktarget}" "${linkname}" 
            elif [[ "$(readlink ${linkname})" != "${linktarget}" ]];then
                eerror "ERROR: Existing path \"${linkname}\" is not a symlink to \"${linktarget}\". Expecting nothing or the symlink there."
            else
                einfo "Link ${linkname} is good."
            fi
        fi
    done

    mkdir -p   "/run/tmp"
    chmod 1777 "/run/tmp"
	populate_run_tree
}
stop(){
	echo -n ""
}

populate_run_tree(){
    if [[ -n "${run_tree}" ]];then
        einfo Populating \"/run\" with tree from \""${run_tree}"\".
        for dir in $(ls -1Ab "${run_tree}");do
            cp -an -t /run "$(realpath -e "${run_tree}/${dir}")"
        done
    fi
}
