# etc/init.d/run_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="run_subdirs"
description=""

depend() {
need root
}

start() {
    einfo "populating run subdirs"
    for node in $(cat /etc/nodelist);do
        hostname="$(hostname)"
        if [[ "${node}" != "${hostname}" ]];then
            mkdir      "/run/${node}"
            chmod 0751 "/run/${node}"
        else
            ln -s "/run" "/run/%{hostname}"
        fi
    done

    mkdir "/run/tmp"
    chmod 1777 "/run/tmp"
}
stop(){
	echo -n ""
}
