# etc/init.d/run_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="run_subdirs"
description="population subdirs and stuff inside of /run"

depend() {
need root
}

start() {
    einfo "populating run subdirs"
    for node in $(cat /etc/nodelist);do
        einfo "node ${node}"

        if [[ "${node}" != "$(hostname)" ]];then
            mkdir -p   "/run/${node}"
            chmod 0751 "/run/${node}"
        else
            # ln -s "." "/run/${node}"
            linkname="/run/${node}"
            linktarget="."
            if ! [[ -e "${linkname}" ]];then
                ln -s  "${linktarget}" "${linkname}" 
            elif [[ "$(readlink ${linkname})" != "${linktarget}" ]];then
                eerror "ERROR: Existing path \"${linkname}\" is not a symlink to \"${linktarget}\". Expecting nothing or the symlink there."
            else
                einfo "Existing link at ${linkname} is valid."
            fi
        fi
    done

    mkdir -p   "/run/tmp"
    chmod 1777 "/run/tmp"
}
stop(){
	echo -n ""
}
