# etc/init.d/run_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="run_subdirs"
description="population subdirs and stuff inside of /run"

depend() {
need root
}

start() {
    einfo "populating run subdirs"
    for node in $(cat /etc/nodelist);do
        einfo "node ${node}"
        if [[ "${node}" != "$(hostname)" ]];then
            mkdir -p   "/run/${node}"
            chmod 0751 "/run/${node}"
        else
            einfo "ln -s /run /run/${node}"
            if ! [[ -e "/run/${node}" ]];then
                ln -s "/run" "/run/${node}"
            elif ! [[ $( readlink  "/run/${node}" ) = "/run" ]];then
                eerror "ERROR: Path \""/run/${node}"\" exists an is not a symlink pointing at \"/run\"."
            fi
        fi
    done

    mkdir -p   "/run/tmp"
    chmod 1777 "/run/tmp"
}
stop(){
	echo -n ""
}
