# etc/init.d/tmp_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="tmp_subdirs"
description=""

depend() {
need localmount run_subdirs tmp-backing
}

start() {
    einfo "populating tmp subdirs"

    for node in $(cat /etc/nodelist);do
        einfo "node ${node}"

        if [[ "${node}" != "$(hostname)" ]];then
            tmpdir="/tmp/${node}"
        else
            tmpdir="/tmp"
            
            linkname="${tmpdir}/${node}"
            linktarget="."
            if ! [[ -e "${linkname}" ]];then
                ln -s  "${linktarget}" "${linkname}" 
            elif [[ "$(readlink ${linkname})" != "${linktarget}" ]];then
                eerror "ERROR: Existing path \"${linkname}\" is not a symlink to \"${linktarget}\". Expecting nothing or the symlink there."
            else
                einfo "Existing link at ${linkname} is valid."
            fi
        fi

        mkdir -p   "${tmpdir}"
        chmod 1777 "${tmpdir}"

        mkdir -p   "${tmpdir}/cache"
        chmod 0751 "${tmpdir}/cache"

        mkdir -p                 "${tmpdir}/cache/ccache"
        chown -R portage:portage "${tmpdir}/cache/ccache"

        mkdir -p                 "${tmpdir}/portage"
        chown -R portage:portage "${tmpdir}/portage"

        mkdir -p                 "${tmpdir}/log"
        chmod 0751               "${tmpdir}/log"

        mkdir -p                 "${tmpdir}/log/portage"
        chown -R portage:portage "${tmpdir}/log/portage"

    done
}
stop(){
	echo -n ""
}
