# etc/init.d/tmp_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="tmp_subdirs"
description=""

depend() {
need localmount run_subdirs
}

start() {
    einfo "populating tmp subdirs"
	mkdir -p /tmp/var/log
    for node in $(cat /etc/nodelist);do

        hostname="$(hostname)"

        if [[ "${node}" != "${hostname}" ]];then
            tmpdir="/tmp/${node}"
        else
            tmpdir="/tmp"
            ln -s "/tmp/${node}" "/tmp"
        fi

        mkdir      "${tmpdir}"
        chmod 1777 "${tmpdir}"

        mkdir      "${tmpdir}/cache"
        chmod 0751 "${tmpdir}/cache"

        mkdir                    "${tmpdir}/cache/ccache"
        chown -R portage:portage "${tmpdir}/cache/ccache"

        mkdir                    "${tmpdir}/portage"
        chown -R portage:portage "${tmpdir}/portage"

        mkdir                    "${tmpdir}/log"
        chmod 0751               "${tmpdir}/log"

        mkdir                    "${tmpdir}/log/portage"
        chown -R portage:portage "${tmpdir}/log/portage"

    done
}
stop(){
	echo -n ""
}
