# etc/init.d/tmp_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="tmp_subdirs"
description=""

depend() {
need localmount
}

start() {
	mkdir -p /tmp/var/log
}
stop(){
	echo -n ""
}
