# etc/init.d/tmp_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="tmp_subdirs"
description=""

depend() {
need localmount run_subdirs
}

start() {
    einfo "populating tmp subdirs"
	mkdir -p /tmp/var/log
    for node in $(cat /etc/nodelist);do

        if [[ "${node}" != "$(hostname)" ]];then
            tmpdir="/tmp/${node}"
        else
            tmpdir="/tmp"
            if [[ -e /tmp ]] && [[ "$(readlink /tmp)" != "/tmp/${node}" ]];then
                eerror "ERROR: Existing path \"/tmp\" is not a symlink to \"/tmp/${node}\". Expecting nothing or the symlink there."
            elif [[ -e /tmp ]];then
                einfo "Existing link at /tmp is valid."
            else
                ln -s "/tmp/${node}" "/tmp" 
            fi
        fi

        mkdir -p   "${tmpdir}"
        chmod 1777 "${tmpdir}"

        mkdir -p   "${tmpdir}/cache"
        chmod 0751 "${tmpdir}/cache"

        mkdir -p                 "${tmpdir}/cache/ccache"
        chown -R portage:portage "${tmpdir}/cache/ccache"

        mkdir -p                 "${tmpdir}/portage"
        chown -R portage:portage "${tmpdir}/portage"

        mkdir -p                 "${tmpdir}/log"
        chmod 0751               "${tmpdir}/log"

        mkdir -p                 "${tmpdir}/log/portage"
        chown -R portage:portage "${tmpdir}/log/portage"

    done
}
stop(){
	echo -n ""
}
