# etc/init.d/tmp_subdirs
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="tmp_subdirs"
description=""

depend() {
need localmount
}

start() {
    einfo "populating tmp subdirs"
	mkdir -p /tmp/var/log
    for node in $(cat /etc/nodelist);do

        mkdir      "/tmp/${node}"
        chmod 1777 "/tmp/${node}"

        mkdir      "/tmp/${node}/cache"
        chmod 0751 "/tmp/${node}/cache"

        mkdir                    "/tmp/${node}/cache/ccache"
        chown -R portage:portage "/tmp/${node}/cache/ccache"

        mkdir                    "/tmp/${node}/portage"
        chown -R portage:portage "/tmp/${node}/portage"

        mkdir                    "/tmp/${node}/log"
        chmod 0751               "/tmp/${node}/log"

        mkdir                    "/tmp/${node}/log/portage"
        chown -R portage:portage "/tmp/${node}/log/portage"

    done
}
stop(){
	echo -n ""
}
